# iface/usr/lib/NetworkManager/dispatcher.d/00-create-natfree-interface
#!/bin/bash

INTERFACE=$1
ACTION=$2

case "$ACTION" in 
	pre-up|up)
		/usr/share/natfree-iface/natfree-iface
		;;
	*)
		;;

esac
