# client/usr/lib/NetworkManager/dispatcher.d/99-update-natfree-hosts
#!/bin/bash


INTERFACE=$1
ACTION=$2

case "$ACTION" in 
	dhcp4-change|dhcp6-change)
		natfree-client UPDATE
		;;
	*)
		;;
esac

